# system_files/usr/lib/hyprbazzite/etc/udev/rules.d/99-lid-action.rules
# Rule to trigger lid action script on lid change
# This rule triggers whenever the lid switch state changes
ACTION=="change", SUBSYSTEM=="button", KERNEL=="lid*", RUN+="/usr/lib/hyprbazzite/etc/hypr/scripts/lidact.sh %s"
